# systemd/warehouse-frontend.service
[Unit]
Description=Product Warehouse Frontend Service (Development Mode)
After=network.target

[Service]
Type=simple
User=user
Group=user
WorkingDirectory=/home/user/Product_Warehouse

# 环境变量
Environment="NODE_ENV=production"

# 启动命令（开发模式）
# 注意：生产环境建议使用静态文件部署，而不是运行开发服务器
ExecStart=/usr/bin/npm run dev -- --host 0.0.0.0 --port 8003

# 重启策略
Restart=always
RestartSec=5
StartLimitInterval=0

# 资源限制
LimitNOFILE=65536

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=warehouse-frontend

[Install]
WantedBy=multi-user.target
